# SYSTEM/opt/penin-monorepo/SYSTEM/etc/udev/rules.d/70-snap.cups.rules
# This file is automatically generated.
# raw-usb
SUBSYSTEM=="tty", ENV{ID_BUS}=="usb", TAG+="snap_cups_cupsd"
# raw-usb
SUBSYSTEM=="usb", TAG+="snap_cups_cupsd"
# raw-usb
SUBSYSTEM=="usbmisc", TAG+="snap_cups_cupsd"
TAG=="snap_cups_cupsd", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_cups_cupsd $devpath $major:$minor"
# raw-usb
SUBSYSTEM=="tty", ENV{ID_BUS}=="usb", TAG+="snap_cups_ippeveprinter"
# raw-usb
SUBSYSTEM=="usb", TAG+="snap_cups_ippeveprinter"
# raw-usb
SUBSYSTEM=="usbmisc", TAG+="snap_cups_ippeveprinter"
TAG=="snap_cups_ippeveprinter", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_cups_ippeveprinter $devpath $major:$minor"
