# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/ssh.service.d/override.conf
[Service]
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
