# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/profile.d/et_manager.sh
alias et-manager="/opt/et_ultimate/venv/bin/python /opt/et_ultimate/et_ultimate_manager.py"
