# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/profile.d/falcon-autonomy.sh
#!/bin/bash
export FALCON_ROOT="1"
export FALCON_AUTONOMOUS="1"
export PYTHONPATH="/root:/root/IAAA_PENIN_OMEGA/core"
export PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root"
