# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/et-evolver.service.d/override.conf
[Service]
CPUAccounting=true
MemoryAccounting=true
IOAccounting=true
# ajuste conforme sua máquina
CPUQuota=200%
MemoryMax=8G
# prazos de parada/repique
TimeoutStopSec=15s
RestartSec=5s
