# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/et_liga_copilotos.service
[Unit]
Description=ETΩ Liga Copilotos - Execução Infinita
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/et_ultimate/agents/copilots/et_liga_copilotos.py --interval 30
Restart=always
RestartSec=1
Environment=PYTHONUNBUFFERED=1

# Opcional: para registrar logs em journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
