# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/penin.service
[Unit]
Description=PENIN Remote Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/penin-daemon
Restart=always
RestartSec=10
StandardOutput=append:/opt/penin/logs/penin.log
StandardError=append:/opt/penin/logs/penin-error.log

[Install]
WantedBy=multi-user.target
