# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/8gpu-sync.service
[Unit]
Description=File Sync to 8GPU Server
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/root/scripts/8gpu-connection/sync_files.sh
StandardOutput=journal
StandardError=journal