# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/penin-sync.service
[Unit]
Description=PENIN Sync Service - GitHub Bidirectional 24/7
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/penin-autosync
ExecStart=/usr/bin/python3 /opt/penin-autosync/start_sync.py
Restart=always
RestartSec=5

# Logs handled by Python script
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target