# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/snap.chromium.daemon.service
[Unit]
# Auto-generated, DO NOT EDIT
Description=Service for snap application chromium.daemon
Requires=snap-chromium-3248.mount
Wants=network.target
After=snap-chromium-3248.mount network.target snapd.apparmor.service
X-Snappy=yes

[Service]
EnvironmentFile=-/etc/environment
ExecStart=/usr/bin/snap run chromium.daemon
SyslogIdentifier=chromium.daemon
Restart=on-failure
RestartSec=3
WorkingDirectory=/var/snap/chromium/3248
TimeoutStopSec=30
Type=simple

[Install]
WantedBy=multi-user.target
