# SYSTEM/opt/penin-monorepo/SYSTEM/etc/profile.d/penin_local_falcon.sh
export USE_LOCAL_FALCON=1
export LOCAL_FALCON_BASE_URL=http://localhost:8010/v1
export LOCAL_FALCON_MODEL_NAME=falcon-mamba-7b
