# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/et_lockdown_firewall.service
[Unit]
Description=Firewall de egress por UID (etbot)
Before=et_brain_operacional.service
Wants=network-pre.target
After=network-pre.target

[Service]
Type=oneshot
ExecStart=/opt/et_ultimate/bin/et_lockdown_firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
