# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/neural_farm.service
[Unit]
Description=IA3 SUPREME Neural Farm 24x7
After=network.target

[Service]
ExecStart=/root/bin/run_neural_farm.sh
Restart=always
RestartSec=2
WorkingDirectory=/root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
