# SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/penin-agent.service
[Unit]
Description=PENIN-Ω-MAX Agent (CLI interactive)
After=network.target local-falcon.service

[Service]
User=root
WorkingDirectory=/root/penin-construtor
EnvironmentFile=/etc/default/penin-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/penin-construtor/penin_go.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
