# SYSTEM/opt/penin-monorepo/SYSTEM/etc/systemd/system/cerebro-unico.service
[Unit]
Description=Cérebro Único - Controlador Central
After=network.target ia-unificada.service
Requires=ia-unificada.service

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/controle_cerebro_unico.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Configurações de recursos
MemoryLimit=1G
CPUQuota=50%

# Configurações de segurança
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target